# etc/service/dovecot/run
#!/bin/sh
exec /usr/sbin/dovecot -F | logger -t dovecot
